# linux_script/proj1-backupcrons.cron
30 15 * * *  ~/Documents/Challenges/linux-script/andre-santos/project1-backup_script_args  ~/Documents